# scripts/ling-resource-monitor.service
[Unit]
Description=lingclaude Resource Monitor Dashboard
After=network.target

[Service]
Type=simple
User=ai
Group=ai
WorkingDirectory=/home/ai/lingclaude/scripts
Environment="PYTHONPATH=/home/ai/lingclaude"
ExecStart=/usr/bin/python3 /home/ai/lingclaude/scripts/resource_monitor_dashboard.py
Restart=always
RestartSec=10
MemoryMax=512M
CPUQuota=100%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
